# _config/systemd/user/rotate.service
[Unit]
Description=Auto rotate screen in 2in1 laptop

[Service]
ExecStart=/home/undg/bin/autorotate auto --daemon --display eDP
